# split.service
[Unit]
Description=Split daemon for i3

[Service]
ExecStart=/usr/local/bin/split.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
